# lib/upstart/querierd.conf
# https://wiki.ubuntu.com/SystemdForUpstartUsers

# querierd - multicast querier for networks without a multicast router 

description	"Multicast querier daemon"

# start in normal runlevels when disks are mounted and networking is available
start on runlevel [2345]

# stop on shutdown/halt, single-user mode and reboot
stop on runlevel [016]

# start daemon
exec /usr/bin/python -m querier.service
